# curator4-1285676296.conf
[program:curator4-1285676296]
command=python3 /home/curator4/curator4-1285676296/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1285676296/test.err.log
stdout_logfile=/home/curator4/curator4-1285676296/test.out.log